# systemd/mtg-monitor.service
[Unit]
Description=WSL Camera and Mic Status Monitor
After=network.target

[Service]
EnvironmentFile=%h/.config/switchbot/api.env
ExecStart=%h/dotfiles/bin/monitor-mtg.sh
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target

